# roles/rps_sync_services/files/rps-sync-scripts.service
[Unit]
Description=RPS sync scripts

[Service]
Type=simple
ExecStart=rps-admin --config-file /etc/rps-admin-tools.json sync
User=rps-sync-scripts
Group=rps-sync-scripts
