# files/group-sync.service
[Unit]
Description=Execute group sync

[Service]
Type=oneshot
ExecStart=/opt/group-sync/bin/group-sync.sh