# files/sync-scripts.service
[Unit]
Description=Execute sync scripts

[Service]
Type=simple
ExecStart=/opt/sync-scripts/bin/group-sync.sh
